# src/wiggle_touch/service/wiggle-touch.service
[Unit]
Description=WiggleTouch service

[Service]
WorkingDirectory=%h
ExecStart=/bin/bash /usr/bin/wiggle-touch-boot.sh

[Install]
WantedBy=default.target